# secretlab-wrt/etc/crontabs/root
*/35 * * * * sh /sbin/cleanup > /dev/null
0 3 * * * /usr/share/openclash/openclash_ipdb.sh
